# transmission/rootfs/etc/s6-overlay/s6-rc.d/transmission/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs Transmission
# ==============================================================================
bashio::log.info "Starting Transmission..."

log_level='--log-debug'

if bashio::config.equals 'log_level' 'info'; then
    log_level='--log-info'
fi
if bashio::config.equals 'log_level' 'error'; then
    log_level='--log-error'
fi


exec /usr/bin/transmission-daemon ${log_level} --foreground --config-dir /config/transmission
